# deploy/oxidize-fip.timer
[Unit]
Description=Periodically reconcile oxidize floating-IP DNAT rules

[Timer]
OnBootSec=30s
OnUnitActiveSec=20s
AccuracySec=5s

[Install]
WantedBy=timers.target
